# razer-battery-colouring.timer
[Timer]
OnBootSec=10
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target

